# packager/scripts/packages/b/bluez/install/install.sh
systemctl enable bluetooth

systemctl enable --global obex
